# systemd/random-sensor.service
[Unit]
Description=Random Sensor Logging Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/random-sensor-service --interval 5 --logfile /tmp/sensor.log --device /dev/urandom
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
